# ops/systemd/jimin-garden-asset-cleanup-worker.timer
[Unit]
Description=Run Jimin Garden asset cleanup worker every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=jimin-garden-asset-cleanup-worker.service

[Install]
WantedBy=timers.target
