# ops/systemd/jimin-garden-webhook-worker.service
[Unit]
Description=Jimin Garden webhook worker trigger
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/srv/jimin-garden/shared/.env
TimeoutStartSec=20
ExecStart=/usr/bin/env bash -lc '/usr/bin/flock -n /run/jimin-garden-webhook-worker.lock /usr/bin/curl --fail --silent --show-error --max-time 15 -X POST "$APP_URL/api/worker/webhook" -H "Authorization: Bearer $CRON_SECRET"'
